# nostr-watchdog.timer
[Unit]
Description=Run nostr-watchdog every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=1s

[Install]
WantedBy=timers.target

